# scripts/etc/systemd/system/qmanager-ttl.service
# /lib/systemd/system/qmanager-ttl.service
[Unit]
Description=QManager TTL/HL Rules Persistence
After=network.target qmanager-setup.service

[Service]
Type=oneshot
RemainAfterExit=yes
ConditionPathExists=/etc/firewall.user.ttl
ExecStartPre=/bin/sleep 5
ExecStart=/bin/sh -c '. /etc/firewall.user.ttl && logger -t qmanager_ttl "Applied TTL/HL settings"'

[Install]
WantedBy=multi-user.target
