# scripts/etc/systemd/system/qmanager-mtu.service
# /lib/systemd/system/qmanager-mtu.service
[Unit]
Description=QManager MTU Settings Persistence
After=network.target qmanager-setup.service

[Service]
Type=simple
ConditionPathExists=/etc/firewall.user.mtu
ConditionFileIsExecutable=/usr/bin/qmanager_mtu_apply
ExecStart=/usr/bin/qmanager_mtu_apply
Restart=no

[Install]
WantedBy=multi-user.target
